# storage/build_tools/docker/config/supervisor/render.conf
[program:render]
 command=/liman/server/storage/liman_render 
 autostart=true 
 autorestart=true 
 stderr_logfile=/liman/logs/render.err.log 
 stdout_logfile=/liman/logs/render.out.log 