# build_tools/DEBIAN/postinst
chmod o+xr /liman/
chmod -R o+xr /liman/sandbox
chmod o+x /liman/extensions

mv /liman/server/nginx.conf /etc/nginx/sites-available/liman.conf
mv /liman/webssh/liman-webssh.service /etc/systemd/system/liman-webssh.service

systemctl enable liman-webssh
systemctl start webssh
systemctl restart nginx
systemctl restart php7.3-fpm
find /liman -type d | grep -v /liman/extensions | while read line ; do chown liman:liman $line; done
php /liman/server/artisan migrate --force

printf "\nKurulum Tamamlandı\n\nYönetici Hesabı oluşturmak yada şifrenizi yenilemek için aşağıdaki komutu çalıştırabilisiniz\n\n\n"
echo "sudo runuser liman -c "php /liman/server/artisan administrator""